# files/systemd/yacs_runner.service
[Unit]
Description=Yacs runner

[Service]
ExecStart={} run-components -y

[Install]
WantedBy=multi-user.target