# files/systemd/yacs.service
[Unit]
Description=Yacs daemon

[Service]
ExecStart={} update-components -y

[Install]
WantedBy=multi-user.target