# files/systemd/yacs_updater.service
[Unit]
Description=Yacs updater

[Service]
ExecStart={} update-components -y

[Install]
WantedBy=multi-user.target
Before=yacs_runner.service